# rpi-app/.debian/postinst.sh
#!/usr/bin/env bash

systemctl enable rpi-app
systemctl start rpi-app
